# mariadb-configs/var/spool/cron/crontabs/btorch
# For more information see the manual pages of crontab(5) and cron(8)
#
# m h  dom mon dow   command
30 */2 * * * /home/btorch/projects/nfe_parsing/nfe_parsing_tosql.py -c /home/btorch/projects/nfe_parsing/config.cfg
